# conf/etc/s6/php-fpm/run
#!/bin/sh

php-fpm --nodaemonize --force-stderr
